# homelab/gv-crm.service
[Unit]
Description=Google Voice voicemail -> Twenty CRM batch
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Edit User and paths to match your home lab.
User=youruser
WorkingDirectory=/home/youruser/gv_crm_pipeline
ExecStart=/home/youruser/gv_crm_pipeline/homelab/run_batch.sh
# Don't let a stuck WhisperX run hang forever.
TimeoutStartSec=1800
